# init/systemd/strongswan.service.in
[Unit]
Description=strongSwan IPsec
After=syslog.target

[Service]
ExecStart=@sbindir@/@ipsec_script@ start --nofork
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
